# Jenkins-aws-go/production/templates/main.service.j2
[Unit]
Description=GO
After=network.target

[Service]
Type=simple
ExecStart={{home_folder}}/main
PIDFile=/var/run/simple

[Install]
WantedBy=default.target
